# supervisord.conf
[supervisord]
nodaemon=true

[program:public_bot]
command=python -m bots.public.main
autorestart=true

[program:admin_bot]
command=python -m bots.admin.main
autorestart=true

[program:webhook]
command=python -m payments.webhook_server
autorestart=true
